# ops/supervisor/sepetak-queue.conf
; =============================================================================
; /etc/supervisor/conf.d/sepetak-queue.conf
; -----------------------------------------------------------------------------
; Laravel queue worker untuk memproses notifikasi email (ShouldQueue):
;   - App\Notifications\MemberRegistrationReceived
;   - App\Notifications\AdminNewMemberRegistered
;   - App\Notifications\AgrarianCaseStatusChanged
;
; Memakai driver `redis` sesuai .env (QUEUE_CONNECTION=redis). Ganti ke
; `database` atau `sync` bila tidak memakai Redis.
;
; Instalasi:
;   sudo cp ops/supervisor/sepetak-queue.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl status sepetak-queue:*
;
; Restart setelah deploy (script deploy.sh sudah memanggilnya otomatis):
;   sudo supervisorctl restart sepetak-queue:*
; =============================================================================

[program:sepetak-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/sepetak.org/current/artisan queue:work redis
    --sleep=3
    --tries=3
    --max-time=3600
    --backoff=30
    --queue=default
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/sepetak-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
